# apps/cron/atrium.cron
0 * * * * wget -O - -q -t 1 http://172.16.0.1/intranet/atrium/cron.php
