# apps/cron/moodle.cron
0 * * * * wget -q -O /dev/null http://172.16.0.1/intranet/moodle/admin/cron.php
